# epv.crontab
* * * * * /usr/bin/python /home/de2hk/emelyanov/Diplom_S_11_v2/main.py[de2hk@XXXXXX XXXXX]$
